# src/debian/sylva-algorithm-runner.cron.d
* * * * *  sylva   flock -n /var/lock/sylva-algorithm-upgrade-run-orders -c "sylva-algorithm-upgrade-run-orders" 1>> /var/log/sylva-algorithm-runner/sylva-algorithm-upgrade-run-orders.log 2>> /var/log/sylva-algorithm-runner/sylva-algorithm-upgrade-run-orders-error.log
*/2 * * * *  sylva   flock -n /var/lock/sylva-algorithm-run-autostart -c "sylva-algorithm-run-autostart" 1>> /var/log/sylva-algorithm-runner/sylva-algorithm-run-autostart.log 2>> /var/log/sylva-algorithm-runner/sylva-algorithm-run-autostart-error.log